# conf/supervisor-program-env.conf.in
[program:${options.program}]
command=${options.command}
process_name=${options.program}
autostart=${options.autostart}
user=${options.user}
directory=${options.directory}
priority=${options.priority}
redirect_stderr=${options.redirect_stderr}
stdout_logfile=${options.stdout_logfile}
environment=${options.environment}
