# conf/supervisor-geotrek-celery.conf.in
[program:geotrek_celery]
command=/opt/geotrek-admin/bin/celery -A geotrek worker -c 1 -n geotrek
process_name=geotrek_celery
autostart=true
user=root
directory=/opt/geotrek-admin
priority=20
redirect_stderr=true
stdout_logfile=/opt/geotrek-admin/var/log/supervisor/geotrek-celery.log
environment = DJANGO_SETTINGS_MODULE=geotrek.settings
