# conf/supervisor.service.in
[Unit]
Description=${parts['django'].project} supervisor
After=network.target

[Service]
Type=simple
ExecStart=${parts['buildout'].directory}/bin/supervisord
ExecReload=${parts['buildout'].directory}/bin/supervisord reload
ExecStop=${parts['buildout'].directory}/bin/supervisord shutdown
User=${parts['buildout'].user}

[Install]
WantedBy=multi-user.target
