# conf/supervisor.conf.in
description     "${parts['django'].project} supervisor"

start on started networking
stop on stopping networking
stop on starting shutdown

setuid ${parts['buildout'].user}
setgid ${parts['buildout'].group}

exec ${parts['buildout'].directory}/bin/supervisord --nodaemon

respawn

# This is required to ensure all started processes to be stopped before Upstart times out
pre-stop script
  ${parts['buildout'].directory}/bin/supervisorctl stop all
end script
