# debian/geotrek.service
[Unit]
Description=Geotrek-admin

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
