# caminae/conf/supervisor.conf.in
description     "supervisor"

start on started network
stop on stopping network
stop on starting shutdown

exec ${parts['buildout'].directory}/bin/supervisord --nodaemon

respawn
