# caminae/conf/supervisor.conf.in
description     "supervisor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec ${parts['buildout'].directory}/bin/supervisord --nodaemon --configuration ${parts['supervisor']['supervisord-conf']}
